# fetcher/com.clementtang.claude-pulse-fetcher-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clementtang.claude-pulse-fetcher-watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--directory</string>
        <string>/Users/clementtang/claude-pulse/fetcher</string>
        <string>python</string>
        <string>-m</string>
        <string>src.health_check</string>
    </array>

    <!-- Run every 2 hours (offset 1h from fetcher) to check if state.json
         last_checked is fresh. Fires a macOS notification if any collector
         hasn't checked in within 5 hours. Local time (GMT+7). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>1</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>3</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>5</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>11</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>13</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>17</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>19</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>23</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/clementtang/.claude/logs/pulse-watchdog-out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/clementtang/.claude/logs/pulse-watchdog-err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/clementtang/claude-pulse/fetcher</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
